# systemd/umami-apprise-notifier.service
[Unit]
Description=Check Umami analytics for recent visitors and notify via Apprise
# Only run if the network is up, since we need to reach the Umami API.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Load credentials and configuration from a protected env file.
# Copy the .env.example to this path and fill in your values.
EnvironmentFile=/etc/umami-apprise-notifier.env

# uv must be on PATH — adjust if installed elsewhere.
# All options are read from environment variables (see EnvironmentFile above),
# so no CLI flags are needed here.  The script picks up UMAMI_API_KEY (or
# UMAMI_USER + UMAMI_PASSWORD as fallback), UMAMI_URL, etc. automatically.
ExecStart=/usr/bin/env uv run %h/.local/bin/umami_apprise_notifier.py

# Harden the service
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# Allow writing state to ~/.local/share/umami-apprise-notifier/
ReadWritePaths=%h/.local/share/umami-apprise-notifier

[Install]
WantedBy=multi-user.target
